# source/node/supervisord.conf
[supervisord]
nodaemon=true


[program:mysqld]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
